# com.somados.updater.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.somados.updater</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOURUSERNAME/Documents/somados/venv/bin/python</string>
        <string>/Users/YOURUSERNAME/Documents/somados/update_flights.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/YOURUSERNAME/Documents/somados</string>
    
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    
    <key>StandardOutPath</key>
    <string>/Users/YOURUSERNAME/Documents/somados/logs/updater.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOURUSERNAME/Documents/somados/logs/updater-error.log</string>
    
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
